# systemd/jwr.service
[Unit]
Description=Rclone backup for jupyter notebook


[Service]
Type=simple
ExecStart=/home/fargus/jwr/start_backup.sh
